# deployment/templates/deployment/supervisord.conf
[supervisord]
logfile={{ projectdir }}/var/log/supervisord.log
pidfile={{ projectdir }}/var/run/supervisord.pid

[unix_http_server]
file={{ projectdir }}/var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://{{ projectdir }}/var/run/supervisord.sock

[program:mysite]
process_name=mysite
command=env DJANGO_SETTINGS_MODULE='settings' {{ projectdir }}/venv/bin/gunicorn -c {{ projectdir }}/etc/gunicorn.conf mysite.wsgi:application
directory={{ projectdir }}
autostart=true
autorestart=true
stderr_logfile={{ projectdir }}/var/log/mysite.err.log
stdout_logfile={{ projectdir }}/var/log/mysite.out.log

{% if varnishd and varnish_port %}
[program:varnish]
command={{varnishd}} -a localhost:{{varnish_port}} -F -n {{projectdir}}/var/lib/varnish/{{sitename}} -f {{projectdir}}/etc/varnish.vcl
stderr_logfile=var/log/varnish.err.log
stdout_logfile=var/log/varnish.out.log
{% endif %}

